# packages/mqttstore/mqttstore-0.4.0.tar.gz/mqttstore-0.4.0/systemd/mqttstore.service
[Unit]
Description=MQTT to SQLite Storage

[Service]
Type=simple
ExecStart=/usr/bin/env python3 -m mqttstore sqlite
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
